# plugins/time-manager/files/etc/uci-defaults/99-fix-racecondition-while-setting-timezone
#!/bin/sh

sed -i 's|echo "$timezone" > /tmp/TZ|echo "$timezone" > /tmp/TZ.tmp \&\& mv /tmp/TZ.tmp /tmp/TZ|g' /etc/init.d/system

exit 0
