# plugins/time-manager/files/etc/uci-defaults/99-remove-default-chronyd-init-script
#!/bin/sh

/etc/init.d/chronyd disable
rm /etc/init.d/chronyd

exit 0
