# plugins/deviceinfo-manager/files/etc/uci-defaults/15_remove_16_prplos_cron_logrotate
#!/bin/sh

# Logrotate is now managed by tr181-deviceinfo

rm -f /etc/uci-defaults/16_prplos_cron_logrotate

exit 0
