# plugins/deviceinfo-manager/files/etc/uci-defaults/99-prevent-uci-setting-hostname
#!/bin/sh

sed -i 's|echo "$hostname" > /proc/sys/kernel/hostname|#echo "$hostname" > /proc/sys/kernel/hostname|g' /etc/init.d/system

exit 0
